# shadowsocks/shadowsocks.ini
[program:shadowsocks]
command=/bin/bash -c \
    '. /ggbandxyz/supervisor/env.profile; \
    [ `docker ps -a -f "name=shadowsocks"| wc -l` -gt 1 ] && \
    exec docker start -a shadowsocks || \
    docker run --name shadowsocks \
        -p 1080:1080 \
        -p 52580:52580 \
        -p 52581:52581 \
        -p 52582:52582 \
        -p 52583:52583 \
        -p 52584:52584 \
        -v /ggbandxyz/shadowsocks/shadowsocks.json:/etc/shadowsocks.json \
        -d \
        shadowsocks:1.0; \
    exec docker start -a shadowsocks'
process_name= shadowsocks
numprocs=1
autostart=true
autorestart=true
startretries = 6
user=root
stdout_logfile=/ggbandxyz/log/shadowsocks.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/ggbandxyz/log/shadowsocks.log
stderr_logfile_maxbytes=10MB
stopwaitsecs=20
